# orchestration/lxc/dev_environment/ctbase/rootfs/etc/systemd/system/runatstartup.service
[Unit]
Description=run script at startup

[Service]
WorkingDirectory=/
ExecStart=/bin/bash runatstartup
Restart=no

[Install]
WantedBy=multi-user.target
